# packages/aa-status-notifications/aa_notifications.service
[Unit]
Description=Look for astroarch news and show desktop notifications
Requires=display-manager.service
After=display-manager.service

[Service]
Type=simple
ExecStart=sleep 90 && /usr/bin/astroarch-status-notifications
Restart=no

[Install]
WantedBy=multi-user.target
