# packages/aa-status-notifications/aa_notifications.service
[Unit]
Description=Look for astroarch news and show desktop notifications

[Service]
User=astronaut
Type=simple
ExecStartPre=sleep 45
ExecStart=/usr/bin/astroarch-status-notifications
Restart=no

[Install]
WantedBy=multi-user.target
