# packages/aa-status-notifications/script.install
post_install () {
  systemctl daemon-reload
  systemctl enable aa_notifications.service
  systemctl start aa_notifications.service
}