# packages/lgpio/lgpio.install
post_install() {
  systemctl daemon-reload
  systemctl enable rgpiod
  systemctl start rgpiod
}
